# tools/osx_packaging/etc/profile.d/libxml2-bin.csh
setenv XML_CATALOG_FILES '/sw/etc/xml/catalog'
setenv SGML_CATALOG_FILES '/sw/etc/sgml/catalog'
