# tools/osx_packaging/etc/profile.d/libxml2-bin.sh
export XML_CATALOG_FILES='/sw/etc/xml/catalog'
export SGML_CATALOG_FILES='/sw/etc/sgml/catalog'
